# tools/overlay_debian12/usr/lib/systemd/system/usb-otg-state.path
[Unit]
Description=Monitor USB OTG State Changes

[Path]
PathModified=/sys/class/udc/49000000.usb-otg/state

[Install]
WantedBy=multi-user.target

